# systemd/kalshi-btc15m-runner.service
[Unit]
Description=Kalshi BTC 15m Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/kalshi-btc-15m-bot
EnvironmentFile=/etc/kalshi-btc15m.env
ExecStart=/root/projects/kalshi-btc-15m-bot/.venv/bin/python -u -m btc15m_bot.runner
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
