# systemd/pm-btc5m-runner.service
[Unit]
Description=Polymarket BTC 5m Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/polymarket-btc-5m-bot
EnvironmentFile=/etc/polymarket-btc5m.env
ExecStart=/root/projects/polymarket-btc-5m-bot/.venv/bin/python -m btc5m_bot.runner
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
